# dot_config/systemd/user/ram-watch.service
[Unit]
Description=RAM watcher with dunst notifications
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/ram-watch.sh

